# scripts/deploy/slingshot.service
[Unit]
Description=Slingshot Trading Engine v4.3.5 Titanium
Documentation=https://github.com/slingshot-trading
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=slingshot
Group=slingshot
WorkingDirectory=/opt/slingshot

# Ejecutar con Uvicorn (Backend FastAPI)
ExecStart=/opt/slingshot/venv/bin/python -m uvicorn engine.api.main_router:app --host 0.0.0.0 --port 8000 --workers 1

# ═══ WATCHDOG: Auto-resurrección ═══
Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=10

# ═══ Aislamiento de Recursos ═══
MemoryMax=2G
MemoryHigh=1.5G
CPUQuota=90%

# ═══ Timeouts ═══
TimeoutStartSec=30
TimeoutStopSec=15

# ═══ Seguridad básica ═══
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/slingshot/tmp /opt/slingshot/logs

# ═══ Variables de entorno ═══
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/slingshot
EnvironmentFile=-/opt/slingshot/.env

# ═══ Logging ═══
StandardOutput=journal
StandardError=journal
SyslogIdentifier=slingshot

[Install]
WantedBy=multi-user.target
